# docker-build/swift/supervisord.conf
[unix_http_server]
file=/var/run/supervisor/supervisor.sock   ; (the path to the socket file)
chmod = 0770
chown = swift:swift
username = dummy
password = dummy

;[inet_http_server]         ; inet (TCP) server disabled by default

[supervisord]
user=swift
logfile=/dev/null
logfile_maxbytes=0
loglevel=info
pidfile=/var/run/supervisor/supervisord.pid
nodaemon=true
minfds=1024
minprocs=200

[supervisorctl]
serverurl=unix:///var/run/supervisor/supervisor.sock
username = dummy
password = dummy

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:memcached]
command=/bin/sh -c "source /etc/swift/memcached && /usr/bin/memcached -v -u $USER -p $PORT -m $CACHESIZE -c $MAXCONN $OPTIONS"
startsecs=3
startretries=1
stopwaitsecs = 3
priority = 1
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:proxy-server]
command=/usr/bin/python /usr/bin/swift-proxy-server -v /etc/swift/proxy-server.conf
startsecs=3
startretries=1
stopwaitsecs = 3
priority = 2
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:rsync]
command=/usr/bin/rsync --config /etc/swift/rsyncd.conf --daemon --no-detach
startsecs=3
startretries=1
stopwaitsecs = 3
priority = 3
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:swift-object-server]
command=/usr/bin/swift-object-server -v /etc/swift/object-server.conf
startsecs=3
startretries=1
stopwaitsecs = 3
priority = 4
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:swift-object-replicator]
command=/usr/bin/swift-object-replicator -v /etc/swift/object-server.conf
startsecs=3
startretries=1
stopwaitsecs = 3
priority = 5
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:swift-object-updater]
command=/usr/bin/swift-object-updater -v /etc/swift/object-server.conf
startsecs=3
startretries=1
stopwaitsecs = 3
priority = 6
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:swift-object-auditor]
command=/usr/bin/swift-object-auditor -v /etc/swift/object-server.conf
startsecs=3
startretries=1
stopwaitsecs = 3
priority = 6
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:swift-container-server]
command=/usr/bin/swift-container-server -v /etc/swift/container-server.conf
startsecs=3
startretries=1
stopwaitsecs = 3
priority = 7
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:swift-container-replicator]
command=/usr/bin/swift-container-replicator -v /etc/swift/container-server.conf
startsecs=3
startretries=1
stopwaitsecs = 3
priority = 8
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:swift-container-updater]
command=/usr/bin/swift-container-updater -v /etc/swift/container-server.conf
startsecs=3
startretries=1
stopwaitsecs = 3
priority = 9
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:swift-container-auditor]
command=/usr/bin/swift-container-auditor -v /etc/swift/container-server.conf
startsecs=3
startretries=1
stopwaitsecs = 3
priority = 10
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:swift-account-server]
command=/usr/bin/swift-account-server -v /etc/swift/account-server.conf
startsecs=3
startretries=1
stopwaitsecs = 3
priority = 11
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:swift-account-replicator]
command=/usr/bin/swift-account-replicator -v /etc/swift/account-server.conf
startsecs=3
startretries=1
stopwaitsecs = 3
priority = 12
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:swift-account-reaper]
command=/usr/bin/swift-account-reaper -v /etc/swift/account-server.conf
startsecs=3
startretries=1
stopwaitsecs = 3
priority = 13
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:swift-account-auditor]
command=/usr/bin/swift-account-auditor -v /etc/swift/account-server.conf
startsecs=3
startretries=1
stopwaitsecs = 3
priority = 14
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -c /etc/swift/nginx.conf
startsecs=3
startretries=1
stopwaitsecs = 3
priority = 15
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0


